# dgx-services/systemd/llama_cpp.service
[Unit]
Description=llama.cpp Service (Kimi-K2 80B)
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/dgx_stack/bin/start_llama_cpp.sh
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
